# redis/redis/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
# 注意这里指定的redis.conf文件路径，必须是绝对路径
command=/software/redis-3.2.11/src/redis-server /redis/redis.conf

